# infrastructure/ansible/panda/files/application.service
[Unit] 
Description=Application

[Service]
ExecStart=/usr/bin/java -jar /home/ubuntu/panda-0.0.1.jar

[Install] 
WantedBy=multi-user.target
